# systemd/mcp-kali.service.in
[Unit]
Description=MCP Kali scheduler and HTTP API
Documentation=https://github.com/dariohy/mcp-kali
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
User=@MCP_KALI_USER@
Group=@MCP_KALI_GROUP@
# Rendered from the selected service account at install time. The service's
# working directory must be writable because tools can create local artifacts.
WorkingDirectory=@MCP_KALI_HOME@
# Keep the application configuration out of systemd's environment: mcp-kali
# reads this file itself, allowing SIGHUP to reload supported settings.
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=@MCP_KALI_BIN@ --config-file @MCP_KALI_CONFIG_FILE@ --system-data-dir @MCP_KALI_SYSTEM_DATA_DIR@ --config-dir @MCP_KALI_CONFIG_DIR@
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s
TimeoutStopSec=15s
UMask=0077
StateDirectory=mcp-kali
ProtectSystem=full
PrivateTmp=true

# Deliberately do not set NoNewPrivileges=yes, PrivateUsers=yes, or a
# restrictive CapabilityBoundingSet=. Root-required declarative tools rely on
# the selected service user's external passwordless sudo policy.

[Install]
WantedBy=multi-user.target
